# conf/app.conf
[program:app]
command = /usr/bin/python /opt/freezing-tyrion/app.py
user=python
group=python 
logfile = /tmp/supervisord.log
